# config/staging/supervisor.conf
[program:med-db]
command = /var/www/med-db/env/bin/gunicorn -w 4 meddb.wsgi:application --bind 0.0.0.0:8080
directory = /var/www/med-db/med-db
user = www-data
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/med-db.log
stderr_logfile = /var/log/supervisor/med-db_err.log

[program:frontend]
command = /var/www/med-db/env/bin/gunicorn -w 4 frontend:app --bind 0.0.0.0:5000
directory = /var/www/med-db
user = www-data
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/med-db.log
stderr_logfile = /var/log/supervisor/med-db_err.log
